# reboottest.service
[Unit]
Description=Reboot test

[Service]
ExecStart=/usr/bin/reboot.sh
KillMode=process
Type=simple

[Install]
WantedBy=multi-user.target
